# autostart/caps-to-escape.desktop
[Desktop Entry]
Type=Application
Name=Caps Lock to Escape
Exec=setxkbmap -option caps:escape
Hidden=true

