# autostart/light-locker.desktop
[Desktop Entry]
Type=Application
Name=Light Locker
Exec=light-locker
Hidden=true
NoDisplay=false
X-GNOME-Autostart-enabled=true

